# packaging/linux/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable portmaster.service

echo "Please reboot your system"